# deploy/hub/neohost-security.service
[Unit]
Description=NeoHost Security Hub (Node.js)
After=network.target mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/neohost-security/hub
Environment=NODE_ENV=production
Environment=SERVE_STATIC=1
Environment=HOST=127.0.0.1
Environment=PORT=7654
Environment=SECURITY_API_TOKEN=schimba-acest-token-secret
Environment=DATABASE_URL=mysql://user:pass@127.0.0.1:3306/neohost
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neohost-security

[Install]
WantedBy=multi-user.target
