# neohost-security.service
[Unit]
Description=NeoHost Security Monitor Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/neohost-security/backend
Environment="SECURITY_API_TOKEN=schimba-acest-token-secret"
Environment="DATABASE_URL=postgresql://user:pass@localhost/neohost"
Environment="HOST=127.0.0.1"
Environment="PORT=7654"
ExecStart=/opt/neohost-security/venv/bin/gunicorn \
    --worker-class flask_sock.gunicorn.Worker \
    --workers 2 \
    --bind 127.0.0.1:7654 \
    --timeout 120 \
    --keep-alive 30 \
    app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neohost-security

[Install]
WantedBy=multi-user.target
